# deerix.cron
*/15 * * * * (cd /config ; /config/deerix_linux)